# infra/systemd/axioma-backup.service
[Unit]
Description=Axioma PostgreSQL Backup Service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/axioma-local/infra/scripts/backup_postgres.sh
